# nexabridge.service
[Unit]
Description=NEXA bridge systemd service
Wants=network-online.target
After=network-online.target
Before=systemd-user-sessions.service

[Service]
TimeoutStartSec=0
WorkingDirectory=/home/bridge/nexabridge
ExecStart=/usr/bin/python run.py
User=bridge
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
